# services/evolution/deploy/evolution/admin-agent/sekailink-postfix-queue-state.service
[Unit]
Description=SekaiLink Postfix Queue Snapshot
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/sekailink/evolution/admin-agent/bin/update-postfix-queue-state.sh \
  /opt/sekailink/evolution/admin-agent/data/postfix_queue_state.json

