# services/evolution/deploy/evolution/api-gateway/room-query/sekailink-evolution-room-query.service
[Unit]
Description=SekaiLink Evolution Room Query Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sekailink
Group=sekailink
WorkingDirectory=/opt/sekailink/evolution/api-gateway
ExecStart=/opt/sekailink/evolution/api-gateway/bin/sekailink_evolution_room_query_service \
  --config /opt/sekailink/evolution/api-gateway/config/evolution_room_query.json
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sekailink/evolution/api-gateway

[Install]
WantedBy=multi-user.target
